# etc/init.d/neatplan
#!/sbin/openrc-run

# neatplan by c0m4r
# https://github.com/c0m4r/neatplan
# Licnse: GPLv3

neatplan_script="/opt/neatplan/neatplan.py"
neatplan_config="/etc/neatplan/default.conf"
neatplan_runfile="/run/neatplan"
ipcmd="/sbin/ip"

depend() {
    need localmount hostname
    want dev-settle
    after bootmisc hwdrivers modules
    provide net
    keyword -jail -prefix -vserver -docker
}

start() {
    # Invoke neatplan script
    $neatplan_script --config $neatplan_config
    # Waiting until tentatives are gone
    eindent
    ifaces=$(cat $neatplan_runfile)
    for iface in $ifaces ; do
        tentatives=$($ipcmd -6 a s $iface | grep tentative)
        iters=1
        until [ ! "$tentatives" ] ; do
            sleep 0.1
            tentatives=$($ipcmd -6 a s $iface | grep tentative)
            iters=$(( iters + 1 ))
            # Give up after 3 seconds to avoid infinite loop
            if [ $iters -gt 30 ]; then
                tentatives=""
                break
            fi
        done
    done
    eoutdent
}
